# systemd/mb-portal.service
[Unit]
Description=MagicBridgeV2 WiFi provisioning captive portal (AP failover)
Documentation=https://github.com/razzrohith/MagicBridgeV2
# Run after the network + wifi have had a chance to come up; the script itself
# waits and only raises the setup AP if there's genuinely no connectivity.
After=network.target wpa_supplicant@wlan0.service systemd-networkd.service
Wants=network.target

[Service]
Type=oneshot
# Run via bash, NOT the script path directly: git checkouts on Windows drop the
# executable bit, so ExecStart=<script> fails with permission-denied and the
# service dies with no log. bash doesn't care about the +x bit.
ExecStart=/bin/bash /opt/magicbridge/provision/mb-portal.sh
# NO start timeout (handoff 26e). This oneshot legitimately blocks waiting for a
# HUMAN to join the hotspot and type WiFi credentials. The old 1200s cap meant a
# user who took >20 min got SIGTERM'd — and because hostapd/dnsmasq live in this
# unit's cgroup, systemd killed the hotspot with it, leaving the unit with no WiFi
# AND no hotspot: exactly the "stranded, power-cycle to recover" failure. Nothing
# is ordered After= this unit, so an indefinite wait does not block boot (verified
# live: during provisioning mb-portal=activating while sshd/kvmd were active).
TimeoutStartSec=infinity
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
